# XiaoRGEEK_PC_For_DS_GFS_TH/ðíR+ã++PC¦µWiFi-ËãÁðí¦Á+·ã¸+-+ÏÍã+Ý+¦V2.5_x1080ÀÍ¦µ--_Release20180706/MulitiVideoTools/video1
#!/bin/sh /etc/rc.common

START=80

start() {
	rmmod gspca_zc3xx > /dev/null 2>&1 && sleep 2
        insmod gspca_zc3xx > /dev/null 2>&1 && sleep 5
	mjpg_streamer -b -i "input_uvc.so -r 320x240 -f 15"  -o "output_http.so -p 8080 -w /web"
	ser2net
}

stop(){
	killall mjpg_streamer
	killall ser2net
}
